# overlay/files-alpine/etc/init.d/coredns
#!/sbin/openrc-run

description="the Coredns DNS server"

pidfile="/run/$RC_SVCNAME.pid"
cfgfile="/etc/$RC_SVCNAME/Corefile"

command="/usr/bin/coredns"
command_args="-quiet -conf=$cfgfile"
command_user="root"
command_background="yes"

depend() {
  need net localmount
  after firewall
}

start_pre() {
  checkpath --file $cfgfile
}