# overlay/files-alpine/etc/init.d/coredns
#!/sbin/openrc-run

description="the Coredns DNS server"

pidfile="/run/$RC_SVCNAME.pid"
cfgfile="/etc/$RC_SVCNAME/Corefile"
logdir="/var/log/$RC_SVCNAME"

command="/usr/bin/coredns"
command_args="$command_args -v 10 -conf=$cfgfile -log_dir=$logdir"
command_user="root"
command_background="yes"

depend() {
  need net localmount tailscaled
  after firewall
}

start_pre() {
  checkpath --file $cfgfile
  mkdir -p $logdir
}