# overlay/files-alpine/etc/init.d/tailscaled
#!/sbin/openrc-run

# supervise the running process
supervisor=supervise-daemon

# require network before starting
depend() {
    need net
}

# run tailscale
name=tailscaled

start_pre() {
    /usr/sbin/tailscaled --cleanup
}

command="/usr/sbin/tailscaled"
command_args="--state=/var/lib/tailscale/tailscaled.state --socket=/run/tailscale/tailscaled.sock >> /var/log/tailscaled.log 2>&1"

stop_post() {
    /usr/sbin/tailscaled --cleanup
}